# crontab
0 */3 * * * /usr/local/bin/php /var/www/html/cron/fetchFeeds.php >> /var/log/lerama.log